# docker/cron/crontab
*/5 * * * * /index_articles.sh >> /var/log/cron.log 2>&1
